# apps/jupyter-ssh-base/config/supervisor/conf.d/jupyter.conf
[program:jupyter]
user=jupyteruser
command=jupyter notebook --config=/home/jupyteruser/.jupyter/jupyter_notebook_config.py
directory=/workspace
autorestart=true
startretries=3
priority=200
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/supervisor/jupyter-stderr.log
stderr_logfile_maxbytes=1048576
redirect_stderr=false
environment=HOME="/home/jupyteruser",USER="jupyteruser",LANG="zh_CN.UTF-8",LC_ALL="zh_CN.UTF-8",VIRTUAL_ENV="/opt/venv",PATH="/opt/venv/bin:/usr/local/bin:/usr/bin:/bin"
