# apps/jupyter-ssh-base/config/supervisor/conf.d/sshd.conf
[program:sshd]
user=root
command=/bin/bash -c "mkdir -p /run/sshd && exec /usr/sbin/sshd -D -e"
directory=/
autorestart=true
startretries=3
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=LANG="zh_CN.UTF-8",LC_ALL="zh_CN.UTF-8"
